# artifacts/d114-pmaports/device-xiaomi-lmi/lmi-splash-release.initd
#!/sbin/openrc-run

name="lmi-splash-release"
description="Release the downstream boot splash before the Phosh session"

command="/usr/sbin/lmi-splash-release"

depend() {
	need lmi-seatd
	before greetd
	after localmount
}
