# artifacts/d114-pmaports/device-xiaomi-lmi/lmi-power-panel.initd
#!/sbin/openrc-run

name="lmi-power-panel"
description="Power-key panel blank and unblank control for lmi"

command="/usr/sbin/lmi-power-panel"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
	need greetd
	after localmount
}
