# artifacts/images/pmos-lmi-v41-downstream-qrtr-ns-cnss-property-full-20260624.manifest
artifact_boot=pmos-lmi-normalboot-v41-downstream-qrtr-ns-cnss-property-20260624.img
artifact_userdata=xiaomi-lmi-v41-downstream-qrtr-ns-cnss-property-userdata-20260624.img
source_boot=/tmp/postmarketOS-export-v41-downstream-qrtr-ns-cnss-property/boot.img
source_userdata=/home/microstar/.local/var/pmbootstrap/chroot_native/home/pmos/rootfs/xiaomi-lmi.img
source_boot_sha256=b39dc93f1583e4d7d2a3e80041a5aa47be57c6fc147efb889aacc926cc62fee8
source_userdata_sha256=e2723bbada061049d83ced8862d0e3e0ed8a7ac0e6a622dbbcdfb62f21d79fc0
kernel_sha256=38c38390ca9a474b4d29d24fb25ad9139bb58e2ad9cd88b5b601abad2f8c2d5e
dtb_sha256=aee89cc172734de955a11ec335b16d3a1b5da51667083b919271c2b6902d57a6
source_ramdisk_sha256=15dccf65815fdacdd2fe56c9a633e13bfe81b7cd42b79ef2b730b93a2bd6db72
v27_ramdisk_sha256=41e3b8504c60b4dc130b1d7c89cd68e6c30ba299ff14d4f42ff16a252d62a0ef
artifact_boot_sha256=e34a87d5f69e8e19e42788440deeb991530fe7846232b748f6e5d008ec07e7ce
artifact_userdata_sha256=e2723bbada061049d83ced8862d0e3e0ed8a7ac0e6a622dbbcdfb62f21d79fc0
artifact_boot_size=52924416
artifact_userdata_size=1819279360
deviceinfo_rootfs_image_sector_size=4096
deviceinfo_usb_network_function=rndis.usb0
deviceinfo_usb_idVendor=0x0525
deviceinfo_usb_idProduct=0xA4A2
loopdevfix=post-export init_functions.sh patch
cmdline=androidboot.hardware=qcom androidboot.console=ttyMSM0 androidboot.memcg=1 lpm_levels.sleep_disabled=1 msm_rtb.filter=0x237 service_locator.enable=1 androidboot.usbcontroller=a600000.dwc3 swiotlb=2048 loop.max_part=7 cgroup.memory=nokmem,nosocket reboot=panic_warm androidboot.fstab_suffix=qcom androidboot.init_fatal_reboot_target=recovery  pmos_boot_uuid=ae5fbef7-e078-44a5-9cfd-862b0369942f pmos_root_uuid=523966fc-c592-4669-a94d-aa7e090f7f25 pmos_rootfsopts=defaults
variant=v41-downstream-qrtr-ns-cnss-property
purpose=downstream boot/rootfs with RNDIS+SSH baseline and manual Wi-Fi bring-up
device_xiaomi_lmi_expected=1-r102
linux_xiaomi_lmi_expected=4.19.325-r8
ssh_strategy=shelli/postmarketos-base OpenRC SSH as in verified v27 baseline
wifi_strategy=SSH-first boot; run static qrtr-ns before CNSS; mount stock Android system/vendor/runtime APEX read-only; create /data/vendor/wifi sockets; run vendor cnss-daemon with LD_PRELOAD Android property shim; manual WLAN trigger
wifi_services_expected=/usr/sbin/lmi-qrtr-ns /etc/init.d/lmi-qrtr-ns /usr/sbin/lmi-wifi-start /etc/conf.d/lmi-wlan-firmware
safety=no flashing or rebooting performed by this build script
