# artifacts/images/pmos-lmi-v45-cnss-status-idempotent-full-20260624.manifest
artifact_boot=pmos-lmi-normalboot-v45-cnss-status-idempotent-20260624.img
artifact_userdata=xiaomi-lmi-v45-cnss-status-idempotent-userdata-20260624.img
source_boot=/tmp/postmarketOS-export-v45-cnss-status-idempotent/boot.img
source_userdata=/home/microstar/.local/var/pmbootstrap/chroot_native/home/pmos/rootfs/xiaomi-lmi.img
source_boot_sha256=cf81a99d35c8e7b60697cab3c37ed075e8d707a524c4c091f7c64000ed000b3b
source_userdata_sha256=18a9111b8e5ce0e76bc01eb47246a6325352cb2686bb9b8954651e706faa9307
kernel_sha256=38c38390ca9a474b4d29d24fb25ad9139bb58e2ad9cd88b5b601abad2f8c2d5e
dtb_sha256=aee89cc172734de955a11ec335b16d3a1b5da51667083b919271c2b6902d57a6
source_ramdisk_sha256=15dccf65815fdacdd2fe56c9a633e13bfe81b7cd42b79ef2b730b93a2bd6db72
v27_ramdisk_sha256=41e3b8504c60b4dc130b1d7c89cd68e6c30ba299ff14d4f42ff16a252d62a0ef
artifact_boot_sha256=e688cad3b23b8635b569b3427eb7aca053ce178804e449d6ed843a9367354b83
artifact_userdata_sha256=18a9111b8e5ce0e76bc01eb47246a6325352cb2686bb9b8954651e706faa9307
artifact_boot_size=52924416
artifact_userdata_size=1819279360
deviceinfo_rootfs_image_sector_size=4096
deviceinfo_usb_network_function=rndis.usb0
deviceinfo_usb_idVendor=0x0525
deviceinfo_usb_idProduct=0xA4A2
loopdevfix=post-export init_functions.sh patch
cmdline=androidboot.hardware=qcom androidboot.console=ttyMSM0 androidboot.memcg=1 lpm_levels.sleep_disabled=1 msm_rtb.filter=0x237 service_locator.enable=1 androidboot.usbcontroller=a600000.dwc3 swiotlb=2048 loop.max_part=7 cgroup.memory=nokmem,nosocket reboot=panic_warm androidboot.fstab_suffix=qcom androidboot.init_fatal_reboot_target=recovery  pmos_boot_uuid=a016f563-8cd3-4d1e-affa-4a27ee26b2aa pmos_root_uuid=9034073b-0ad5-45e3-a8f0-59e7266bdb92 pmos_rootfsopts=defaults
variant=v45-cnss-status-idempotent
purpose=downstream boot/rootfs with RNDIS+SSH baseline and manual Wi-Fi bring-up
device_xiaomi_lmi_expected=1-r106
linux_xiaomi_lmi_expected=4.19.325-r8
ssh_strategy=shelli/postmarketos-base OpenRC SSH as in verified v27 baseline
wifi_strategy=SSH-first boot; run static qrtr-ns before CNSS; mount stock Android system/vendor/runtime APEX read-only; mount persist read-only for WLAN MAC; create /data/vendor/wifi sockets; link vendor wlan/qca_cld/WCNSS_qcom_cfg.ini and wlan_mac.bin into /lib/firmware; run vendor cnss-daemon with LD_PRELOAD Android property shim; manual WLAN trigger
wifi_services_expected=/usr/sbin/lmi-qrtr-ns /etc/init.d/lmi-qrtr-ns /usr/sbin/lmi-wifi-start /etc/conf.d/lmi-wlan-firmware
safety=no flashing or rebooting performed by this build script
